# crontab
14 14 * 2 5 /bin/bash ~/modul1/soal1.sh
14 14 14 2 * /bin/bash ~/modul1/soal1.sh
0 * * * * /bin/bash ~/modul1/soal4.sh
2-30/6 * * * * /bin/bash ~/modul1/soal5.sh
